# supervisor/laravel.ini
[supervisord]
nodaemon=true

[program:worker]
process_name=%(program_name)s_%(process_num)02d
command=/var/www/html/artisan queue:work --timeout=60 --tries=5
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/supervisor.log

[program:cron]
command=/usr/sbin/crond -f
autostart=true
autorestart=true

[program:php-fpm]
command=docker-php-entrypoint php-fpm
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true

[program:openresty]
command=/usr/local/openresty/bin/openresty -g 'daemon off;'
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
